# vat-frontend.service
[Unit]
Description=VAT Next.js Frontend Application
After=network.target
Wants=vat-fastapi.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vat
Environment=NODE_ENV=production
Environment=NEXT_PUBLIC_API_BASE_URL=https://www.vertexassistant.ru/api
Environment=NEXT_PUBLIC_TELEGRAM_BOT_NAME=VertexAIassistantBOT
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Ограничения ресурсов
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
